# systemd/.config/systemd/user/battery-notify.service
[Unit]
Description=Battery Notification Service

[Service]
Type=oneshot
ExecStart=/home/alterwill/.local/bin/battery-notify.sh

